# snappymail/rootfs/etc/s6-overlay/s6-rc.d/snappymail-init/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# SnappyMail initialization
# ==============================================================================

DATA_DIR="/var/www/snappymail/data"
CONFIG_DIR="${DATA_DIR}/_data_/_default_"
ADMIN_PASS_FILE="${CONFIG_DIR}/admin_password.txt"

bashio::log.info "Initializing SnappyMail..."

# Create data directory if not exists
mkdir -p "${CONFIG_DIR}"
chown -R snappymail:snappymail "${DATA_DIR}"

# Set admin password from options
ADMIN_PASSWORD=$(bashio::config 'admin_password')
if bashio::var.has_value "${ADMIN_PASSWORD}"; then
    bashio::log.info "Setting admin password..."
    echo "${ADMIN_PASSWORD}" > "${ADMIN_PASS_FILE}"
    chown snappymail:snappymail "${ADMIN_PASS_FILE}"
    chmod 600 "${ADMIN_PASS_FILE}"
fi

# Pre-configure local mailserver if running
LOCAL_DOMAIN="biker633.ddnss.de"
DOMAINS_CONFIG="${CONFIG_DIR}/domains"
mkdir -p "${DOMAINS_CONFIG}"

# Create default domain config for local mailserver
if [ ! -f "${DOMAINS_CONFIG}/${LOCAL_DOMAIN}.ini" ]; then
    bashio::log.info "Creating default domain config for ${LOCAL_DOMAIN}..."
    cat > "${DOMAINS_CONFIG}/${LOCAL_DOMAIN}.ini" << EOF
imap_host = "${LOCAL_DOMAIN}"
imap_port = 993
imap_secure = "SSL"
smtp_host = "${LOCAL_DOMAIN}"
smtp_port = 587
smtp_secure = "TLS"
EOF
    chown snappymail:snappymail "${DOMAINS_CONFIG}/${LOCAL_DOMAIN}.ini"
fi

# Create Gmail domain config
if [ ! -f "${DOMAINS_CONFIG}/gmail.com.ini" ]; then
    bashio::log.info "Creating Gmail domain config..."
    cat > "${DOMAINS_CONFIG}/gmail.com.ini" << EOF
imap_host = "imap.gmail.com"
imap_port = 993
imap_secure = "SSL"
smtp_host = "smtp.gmail.com"
smtp_port = 587
smtp_secure = "TLS"
EOF
    chown snappymail:snappymail "${DOMAINS_CONFIG}/gmail.com.ini"
fi

bashio::log.info "SnappyMail initialization complete!"
