# snappymail/rootfs/etc/s6-overlay/s6-rc.d/snappymail-init/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# SnappyMail initialization
# ==============================================================================

DATA_DIR="/var/www/snappymail/data"
CONFIG_DIR="${DATA_DIR}/_data_/_default_"

bashio::log.info "Initializing SnappyMail..."

mkdir -p "${CONFIG_DIR}/domains"
chown -R snappymail:snappymail "${DATA_DIR}"

ADMIN_PASSWORD=$(bashio::config 'admin_password')
echo "${ADMIN_PASSWORD}" > "${CONFIG_DIR}/admin_password.txt"
chmod 600 "${CONFIG_DIR}/admin_password.txt"

bashio::log.info "SnappyMail initialization complete!"

# Keep process alive (required for longrun)
exec tail -f /dev/null
