# snappymail/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bashio

WEB_PORT=$(bashio::config 'web_port')
sed -i "s/listen [0-9]* default_server/listen ${WEB_PORT} default_server/" /etc/nginx/http.d/snappymail.conf

UPLOAD_MAX_SIZE=$(bashio::config 'upload_max_size')
sed -i "s/client_max_body_size .*/client_max_body_size ${UPLOAD_MAX_SIZE};/" /etc/nginx/http.d/snappymail.conf

bashio::log.info "Starting nginx on port ${WEB_PORT}..."
exec nginx -g "daemon off;"
