# snappymail/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start PHP-FPM
# ==============================================================================

MEMORY_LIMIT=$(bashio::config 'memory_limit')

# Update PHP memory limit
sed -i "s/memory_limit = .*/memory_limit = ${MEMORY_LIMIT}/" /etc/php84/php.ini

exec php-fpm84 -F -R
